# proyecto-si783-2026-i-u1-migrador-de-bd-codigo-web/despliegue/supervisor_migrador-bd.conf
[program:migrador-bd]
; Configuración de Supervisor para Gunicorn
; Lugar: /etc/supervisor/conf.d/migrador-bd.conf

user=www-data
directory=/opt/migrador-bd

; Gunicorn con Flask-SocketIO (usar worker threads en Python 3.12)
command=/opt/migrador-bd/venv/bin/gunicorn --workers 2 --worker-class gthread --threads 4 --bind 127.0.0.1:8000 --timeout 300 --access-logfile /var/log/migrador-bd-access.log --error-logfile /var/log/migrador-bd-error.log wsgi:app

; Autostart y reinicio automático
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10

; Logs
stderr_logfile=/var/log/migrador-bd-error.log
stdout_logfile=/var/log/migrador-bd-access.log

; Variables de entorno
environment=PATH="/opt/migrador-bd/venv/bin",PYTHONUNBUFFERED=1,FLASK_ENV=production,SOCKETIO_ASYNC_MODE=threading

; Límites de recursos
numprocs=1
